# rest/supervisor.conf
[program:restdemo]
command=python3 /home/cupen/rest/main.py
directory=/home/cupen/rest/
autorestart=true
user=cupen
